# templates/pacoloco.service
[Unit]
Requires=unraid-public.mount
After=unraid-public.mount
